# cron/nightly.cron
SHELL=/bin/bash
0 1 * * * root module load python && source /bi/home/andrewss/capstone_monitor/venv/bin/activate && /bi/home/andrewss/capstone_monitor/cron/nightly_stats_collection.py > /dev/null